# supervisor/django.conf
[program:acapp]
command=daphne -b 0.0.0.0 -p 5015 acapp.asgi:application
directory=/root/acapp
autorestart=true
autostart=true
user=root
redirect_stderr=true
stdout_logfile_maxbytes=1MB
stdout_logfile=/root/acapp/daphne.log

